# python/riptide/proxy/riptide@.service
[Unit]
Description=Riptide

[Service]
ExecStart=/usr/bin/riptide_proxy --user=%i
Restart=on-failure

[Install]
WantedBy=multi-user.target
